# puppet/modules/invenio/files/supervisor.conf.d/celeryd.ini
[program:celeryd]
command = /home/vagrant/virtualenv_cmd.sh /home/vagrant/.virtualenvs/pu celeryd -E -A invenio.celery.celery --loglevel=INFO --workdir=/home/vagrant/.virtualenvs/pu
