# puppet/modules/invenio/files/supervisor.conf.d/flower.ini
[program:flower]
command = /home/vagrant/virtualenv_cmd.sh /home/vagrant/.virtualenvs/pu flower --broker=redis://localhost:6371/1
